# files/etc/init.d/auto-modem
#!/bin/sh /etc/rc.common
# 不使用modem，取消这个文件的内容。
#START=99
#STOP=10

#start() {
#    sleep 30
#    ifup modem
#}

#stop() {
#    ifdown modem
#}
